# init.d/obot.service
[Unit]
Description=OBOT - 24/7 channel daemon
After=network-online.target
Wants=network-online.target

[Service]
User=obot
Group=obot
ExecStart=/usr/local/bin/obot -n

[Install]
WantedBy=multi-user.target
