# init.d/obot.service
[Unit]
Description=obot

[Service]
ExecStart=/usr/local/bin/obot -k -p -d /var/lib/obot --logdir /var/log/obot -l debug

[Install]
WantedBy=multi-user.target
