# init.d/obot.service
[Unit]
Description=pure python3 channel bot

[Service]
ExecStart=/usr/local/bin/obot

[Install]
WantedBy=multi-user.target
